# SmartBus/smartbusFeedback.service
[Unit]
Description= Listen smartbus Feedback

[Service]
ExecStart=/mnt/Smartbus/smartbusFeedback.sh

[Install]
WantedBy=multi-user.target
